# server1/wireguard/wg-split-routing-update.timer
[Unit]
Description=Daily update of RU ipset

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
